# supervisor/elasticsearch.conf
[program:elasticsearch]
directory=/
command=/usr/share/elasticsearch/bin/elasticsearch
autostart=true
autorestart=true