# supervisor/logstash.conf
[program:logstash]
directory=/opt/logstash
command=/opt/logstash/bin/logstash agent -f /etc/logstash/conf.d
autostart=true
autorestart=true
